# dotfiles/home.arch/.config/systemd/user/irc-notifications.service
[Unit]
Description=Dunst Daemon

[Service]
ExecStart=/home/markl/GitApps/IRCNotify/client.arch.js

[Install]
WantedBy=default.target

